# etc/systemd/system/patissier.service
[Unit]
Description=patissier service

[Service]
ExecStart=/home/pi/bin/event/patissier.sh
Type=simple

[Install]
WantedBy=multi-user.target